# systemd/on-failure/services/OneShot.service
[Unit]
Description=Example OneShot Service
OnFailure=OnFailure.service

[Service]
Type=oneshot
ExecStart=/etc/systemd/user/OneShot.sh
