# systemd/on-failure/services/OnFailure.service
[Unit]
Description=Example OnFailure Service

[Service]
Type=oneshot
ExecStart=/etc/systemd/user/OnFailure.sh
